# deploy/postgres-mcp.service
[Unit]
Description=PostgreSQL MCP Server
Documentation=https://github.com/UshurbakiyevDavlat/postgres-mcp-server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/PostgresMCP
EnvironmentFile=/opt/PostgresMCP/.env
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=postgres-mcp

[Install]
WantedBy=multi-user.target
